# templates/systemd.mount.j2
[Unit]
Description=Mount {{ name }}

[Mount]
What={{ source }}
Where={{ destination }}

[Install]
WantedBy={{ target }}