# templates/systemd.automount.j2
[Unit]
Description=Automount {{ name }}

[Automount]
Where={{ destination }}

[Install]
WantedBy={{ target }}